# ops/live-vps/runtime/systemd/darkmesh-backup.timer
[Unit]
Description=Run darkmesh backup daily

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
