# ops/live-vps/runtime/systemd/darkmesh-config-prune.service
[Unit]
Description=DarkMesh config backup retention prune
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/darkmesh-config-backup-prune.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
