# ops/live-vps/runtime/systemd/darkmesh-config-prune.timer
[Unit]
Description=Run DarkMesh config backup prune weekly

[Timer]
OnCalendar=Sun *-*-* 04:40:00
Persistent=true
Unit=darkmesh-config-prune.service

[Install]
WantedBy=timers.target
